# deploy/amble-weekly.timer
[Unit]
Description=Run Amble data refresh weekly

[Timer]
OnCalendar=Mon *-*-* 03:15:00 Asia/Singapore
Persistent=true
RandomizedDelaySec=10m
Unit=amble-weekly.service

[Install]
WantedBy=timers.target
