# deploy/amble-weekly.service
[Unit]
Description=Amble weekly approved-data refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=amble
Group=amble
WorkingDirectory=/srv/amble
EnvironmentFile=-/etc/amble.env
ExecStart=/usr/bin/npm run weekly:refresh
Nice=10
NoNewPrivileges=true
PrivateTmp=true
