# deploy/whats-here-weekly.cron
# Replace /srv/whats-here before installing with crontab.
CRON_TZ=Asia/Singapore
15 3 * * 1 cd /srv/whats-here && /usr/bin/npm run weekly:refresh >> /var/log/whats-here-weekly.log 2>&1
